# scripts/com.notion-obsidian-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.notion-obsidian-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>source ~/.zshrc && node /Users/kimuratakuya/line-harness/scripts/sync-notion-to-obsidian.mjs all && qmd embed</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/opt/homebrew/bin:/usr/local/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/kimuratakuya/line-harness</string>
    <key>StandardOutPath</key>
    <string>/Users/kimuratakuya/Library/Logs/notion-obsidian-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kimuratakuya/Library/Logs/notion-obsidian-sync.log</string>
</dict>
</plist>
